# config/clock-weather-xinit.service
[Unit]
Description=Clock and Weather Display on PiTFT
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Environment=FRAMEBUFFER=/dev/fb0
Environment=HOME=/home/pi
WorkingDirectory=/home/pi
ExecStart=/usr/bin/xinit /home/pi/clock_weather.py -- :0 vt1
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
